# userinfo-parent/userinfo-script-support/57-zj-hive/crontab
00 02 * * * sh /data/bash/user_info/zj_hive.sh > /data/bash/user_info/zj_hive.log 2>&1 &